# roles/base/files/profile.d/base.sh
python2 /opt/bin/sysinfo.py

export HISTTIMEFORMAT="%Y-%m-%d %H:%M:%S "
HISTTIMEFORMAT="%Y-%m-%d %H:%M:%S ${USER_IP} `whoami` "
